# config/akilimo/supervisor/conf.d/nginx.conf
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
stdout_logfile=/var/log/supervisor/nginx-out.log
redirect_stderr=true
priority=20  # Start after PHP-FPM
inherit=defaults