# docker/supervisor/common.conf
[program:defaults]
process_name=%(program_name)s_%(process_num)02d
user=www-data
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
redirect_stderr=true
startsecs=5
startretries=3
stopwaitsecs=10