# docker/supervisor/conf.d/laravel-opt-out-queue.conf
[program:laravel-opt-out]
process_name=%(program_name)s_%(process_num)02d
directory=/var/www/html/fuelrod
command=php artisan queue:work --queue=opt-out --sleep=%(ENV_QUEUE_SLEEP)s --tries=%(ENV_QUEUE_MAX_TRIES)s
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
numprocs=1
stdout_logfile=/var/log/supervisor/laravel-opt-out.log
stderr_logfile=/var/log/supervisor/laravel-opt-out-err.log
stopwaitsecs=%(ENV_QUEUE_TIMEOUT)s
priority=400
inherit=defaults